# supervisord/masquerade.conf
[program:masquerade]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/masquerade
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/masquerade.log